# mouse/81-mx-master.rules
ACTION=="add", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c52b",  RUN+="/home/manu/.config/manu/mouse-udev-receiver.sh"
